# crontab
* * * * * /fluentd/etc/downloader.py >> /var/log/cron.log 2>&1

